# etc/systemd/system/letsencrypt.timer
[Unit]
Description=Monthly letsencrypt certificate renewal

[Timer]
OnUnitActiveSec=30d
AccuracySec=1d

[Install]
WantedBy=multi-user.target
